# supervisor/file-upload-server.conf
[program:file-upload-server]
user=ubuntu
command=/home/ubuntu/.virtualenvs/file-upload-server/bin/uwsgi --ini /data/web/file-upload-server/uwsgi.ini
autostart=true
